# playbooks/templates/container_sharry.service.j2
[Unit]
Description={{ container_name }} Container Service

[Service]
Restart=on-failure
RestartSec=30

ExecStartPre=/usr/bin/rm -f /%t/%n-pid /%t/%n-cid
ExecStart=/usr/bin/podman run --conmon-pidfile /%t/%n-pid --cidfile /%t/%n-cid -d -v /var/lib/hodgins/{{ container_name }}/conf:/data/sharry/conf:z -v /var/lib/hodgins/{{ container_name }}/db:/tmp:z -e BASE_URL={{ base_url }} -p {{ container_port }}:9090 --name service_{{ container_name }} {{ image_name }}
ExecStop=/usr/bin/sh -c "/usr/bin/podman rm -f `cat /%t/%n-cid`"

KillMode=none
Type=forking
PIDFile=/%t/%n-pid

[Install]
WantedBy=multi-user.target
